# dot_config/systemd/user/agent-dashboard-refresh.timer
[Unit]
Description=Refresh Agent Dashboard every 30 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
