# dot_config/systemd/user/agent-dashboard-refresh.service
[Unit]
Description=Refresh Agent Dashboard HTML

[Service]
Type=oneshot
ExecStart=%h/.local/bin/generate-dashboard
